# etc/systemd/system/letsencrypt@.service
[Unit]
Description=letsencrypt generate a new certificate %I
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm --name certbot_%i -v "/apps/letsencrypt/etc:/etc/letsencrypt" certbot/certbot:latest -c "/etc/letsencrypt/config/%i.ini" certonly

[Install]
WantedBy=letsencrypt.target

